# tool/image/test/cryptsetup/test_base/usr/lib/systemd/system/initrd-cryptsetup.path
# /usr/lib/systemd/system/initrd-cryptsetup.path
# This file is part of https://github.com/random-archer/mkinitcpio-systemd-tool

[Unit]
Description=Initrd Cryptsetup Path
Documentation=https://github.com/random-archer/mkinitcpio-systemd-tool/blob/master/README.md
ConditionPathExists=/etc/crypttab
ConditionPathExists=/etc/initrd-release
DefaultDependencies=no
Conflicts=shutdown.target
Before=paths.target shutdown.target
# require explicitly to include in initramfs image
Requires=initrd-cryptsetup.service

[Path]
MakeDirectory=yes
DirectoryNotEmpty=/run/systemd/ask-password

[Install]
WantedBy=sysinit.target

[X-SystemdTool]

# disable default password agents
InitrdPath=/usr/lib/systemd/system/systemd-ask-password-console.path replace=yes create=yes
InitrdPath=/usr/lib/systemd/system/systemd-ask-password-wall.path    replace=yes create=yes

# /etc/systemd/system/initrd-cryptsetup.path.d/override.conf

[X-SystemdTool]

# disable default password agents, AGAIN
InitrdPath=/usr/lib/systemd/system/systemd-ask-password-console.path replace=yes create=yes
InitrdPath=/usr/lib/systemd/system/systemd-ask-password-wall.path    replace=yes create=yes
