# tool/image/test/cryptsetup/test_base/usr/lib/systemd/system/initrd-cryptsetup.path
# /usr/lib/systemd/system/initrd-cryptsetup.path
# This file is part of https://github.com/random-archer/mkinitcpio-systemd-tool

# Configure crypto disks in initramfs

# Provides replacement for mkinitcpio hook "sd-encrypt" included in cryptsetup
# https://www.archlinux.org/packages/core/x86_64/cryptsetup/
# https://wiki.archlinux.org/index.php/dm-crypt/System_configuration#Using_sd-encrypt_hook

# note:
# this is a twin unit for initrd-cryptsetup.service
# enable only initrd-cryptsetup.path, initrd-cryptsetup.service is activated on demand

# note:
# remove `sd-encrypt` from HOOKS=()

# note:
# mutually exclusive with initrd-plymouth.path

# service dependencies:
# https://www.archlinux.org/packages/core/x86_64/cryptsetup/

[Unit]
Description=Initrd Cryptsetup Path
Documentation=https://github.com/random-archer/mkinitcpio-systemd-tool/blob/master/README.md
ConditionPathExists=/etc/crypttab
ConditionPathExists=/etc/fstab
ConditionPathExists=/etc/initrd-release
DefaultDependencies=no
Conflicts=shutdown.target
Before=paths.target shutdown.target cryptsetup.target

[Path]
# bind to the twin unit
Unit=initrd-cryptsetup.service
# create monitored folder when missing
MakeDirectory=yes
# activate twin unit when cryptsetup request files are present
DirectoryNotEmpty=/run/systemd/ask-password

[Install]
WantedBy=sysinit.target

[X-SystemdTool]

# provision disk tables in initramfs
InitrdPath=/etc/crypttab   source=/etc/mkinitcpio-systemd-tool/config/crypttab   replace=yes
InitrdPath=/etc/fstab      source=/etc/mkinitcpio-systemd-tool/config/fstab      replace=yes

# /etc/systemd/system/initrd-cryptsetup.path.d/override.conf

[X-SystemdTool]

# inject busybox enviro
InitrdPath=/etc/profile source=/etc/systemd/system/profile.sh
