# tool/image/test/dropbear/etc/systemd/system/initrd-dropbear.service.d/override.conf

# enable password logins

[Service] 
ExecStart=
ExecStart=/bin/dropbear -j -k -F

[X-SystemdTool]
InitrdBuild=/usr/lib/mkinitcpio-systemd-tool/initrd-build.sh command=do_root_login_enable
