# tool/image/test/dropbear/etc/systemd/system/initrd-dropbear.service.d/override.conf

[Service]

# change port
ExecStart=
ExecStart=/bin/dropbear -s -j -k -F -p 221
