# tool/image/test/cryptsetup/etc/systemd/system/initrd-cryptsetup.service.d/override.conf

[Service]

# change timeout
RestartSec=5s
