# tool/image/test/tinysshd/etc/systemd/system/initrd-tinysshd.service.d/override.conf

[Service]

# change port
ExecStart=
ExecStart=/usr/bin/busybox tcpsvd -v 0 222 /usr/bin/tinysshd -v /etc/tinyssh/sshkeydir
