# tool/image/test/tinysshd/etc/systemd/system/initrd-tinysshd.service.d/override.conf

[Service]

# change port number

Environment=SSHD_PORT=222
